# systemd/shadow_validator.service
[Unit]
Description=Hydroponic Shadow Validation Service
After=network.target

[Service]
Type=oneshot
User={{USER}}
Group={{USER}}
WorkingDirectory={{PROJECT_DIR}}
Environment=PYTHONPATH={{PROJECT_DIR}}
ExecStart={{PYTHON_PATH}} {{PROJECT_DIR}}/scripts/shadow_validator.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hydro-shadow-validator

# Resource limits
MemoryMax=512M
CPUQuota=75%
TimeoutStartSec=3600

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{PROJECT_DIR}} /home/{{USER}}/hydro /var/log/hydro